# switch/intel/disable_gpu.service
[Unit]
Description=Disable GPU
 
[Service]
ExecStart=/etc/switch/disable_gpu.sh

[Install]
WantedBy=multi-user.target
